# mrt.init
#!/bin/bash
#
# mrtd:         Starts the Dynamic Route Daemon 
#
# chkconfig: 345 80 45
# description: Dynamic Route Daemon for IPv6 routers
#
# processname: mrtd
# config: /etc/mrtd.conf
#
export PATH=/bin:/sbin:/usr/bin:/usr/sbin

# Sanity checks.
[ -f /etc/mrtd.conf ] || exit 0
[ -x /usr/sbin/mrtd ] || exit 0

# Source function library.
. /etc/rc.d/init.d/functions

# Check that networking is up.
[ "${NETWORKING}" = "no" ] && exit 0

# See how we were called.
case "$1" in
  start)
	show "Starting Route Daemon: "
	daemon +1 mrtd 
	touch /var/lock/subsys/mrtd
	;;
  stop)
	show "Stopping Route Daemon: "
	killproc mrtd 
	rm -f /var/lock/subsys/mrtd
	;;
  status)
	status mrtd
	;;
  restart)
	$0 stop
	$0 start
	;;
  *)
	echo "Usage: $0 {start|stop|status|restart}"
	exit 1
esac

exit 0
